# batch/cron/crontab
0 22 20-30 11 * /srv/www/suilabs/BHV-Marais/batch/cron/mysq-queries.sh
0 22 1-24 12 * /srv/www/suilabs/BHV-Marais/batch/cron/mysq-queries.sh
